# codebase/scorpius/ansible/roles/common/files/etc-environment.target
[Unit]
Description=/etc/environment available and hostname set
Requires=etc-environment.service
After=etc-environment.service
