# codebase/scorpius/ansible/roles/services/files/dcos-journalctl-filebeat.target
[Unit]
Description=DCOS journalctl parser to filebeat
Requires=dcos-journalctl-filebeat.service
After=dcos-journalctl-filebeat.service
